# units/ax25ports.target
[Unit]
Description=All ax.25 ports are up
