# units/radio.target
[Unit]
Description=ax.25 services for Linux

[Install]
WantedBy=default.target
